# private_dot_config/private_systemd/private_user/private_gammastep.service.d/private_override.conf
[Unit]
PartOf=
PartOf=wayland-session.target
After=
After=wayland-session.target

[Install]
WantedBy=
WantedBy=wayland-session.target
