# private_dot_config/private_systemd/private_user/private_mbsync@gmail.service.d/private_override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/proxychains /usr/bin/mbsync -c %h/.config/isync/mbsyncrc -V %i