# private_dot_config/private_systemd/private_user/private_xdg-desktop-portal-wlr.service.d/private_override.conf
[Unit]
PartOf=
PartOf=sway.service
After=
After=sway.service