# private_dot_config/private_systemd/private_user/private_udiskie.service
[Unit]
Description=Removable disk automounter
Documentation=man:udiskie(8)

[Service]
Type=simple
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=multi-user.target
