# private_dot_config/private_systemd/private_user/private_mconnect.service
[Unit]
Description=KDE Connect protocol implementation in Vala/C

[Service]
ExecStart=/usr/bin/mconnect -d

[Install]
WantedBy=multi-user.target
