# private_dot_config/private_systemd/private_user/private_vdirsyncer@.timer
[Unit]
Description=Calendar synchronization timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=15m
Unit=vdirsyncer@%i.service

[Install]
WantedBy=timers.target
