# private_dot_config/private_systemd/private_user/private_vdirsyncer@google.service.d/private_override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/proxychains /usr/bin/vdirsyncer sync %i