# private_dot_config/private_systemd/private_user/private_goimapnotify@gmail.service.d/private_override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/proxychains /usr/bin/goimapnotify -conf %h/.config/imapnotify/%i.conf