# private_dot_config/private_systemd/private_user/private_waybar.service.d/private_override.conf
[Unit]
PartOf=
After=
PartOf=wayland-session.target
After=wayland-session.target

[Service]
Restart=no

[Install]
WantedBy=
WantedBy=wayland-session.target
