# private_dot_config/private_systemd/private_user/private_xdg-desktop-portal.service.d/private_override.conf
[Unit]
After=sway.service
Requires=sway.service